# controller/etc/cron.d/tunnel
@reboot autossh -M 0 -f fas@node.mide.dk -R 2224:localhost:22 -N -o "ServerAliveInterval 45" -o "ServerAliveCountMax 2"
